# etc/tagion-ops-mngr.service
[Unit]
Description=Service for managing tagion operational tests

[Service]
ExecStart=/path/to/myscript.sh
