# etc/tagion-ops-mngr.service
[Unit]
Description=Service for managing tagion operational tests

[Service]
ExecStart=%h/tagion/ops/activate.sh
WorkingDirectory=%h/tagion/ops/
