# src/bin-wave/neuewelle.service
[Unit]
Description="Test deployment of tagionwave"

[Service]
WorkingDirectory=%h/.local/share/tagion/wave
ExecStart=%h/.local/share/bin/tagion neuewelle --keys %h/.local/share/tagion/wallets -p wallet1:0001 -p wallet2:0002 -p wallet3:0003 -p wallet4:0004 -p wallet5:0005
Restart=no
RestartSec=30

[Install]
WantedBy=default.target
