# src/bin-wave/neuewelle.service
[Unit]
Description="Test deployment of tagionwave"

[Service]
WorkingDirectory=%h/.local/share/tagion/wave
ExecStart=%h/.local/share/bin/tagion neuewelle
Restart=no
RestartSec=30

[Install]
WantedBy=default.target
