# src/bin-tagionshell/tagionshell.service
[Unit]
Description="Test deployment of tagionshell"
Requires=network-online.target

[Service]
WorkingDirectory=%h/.local/share/tagion/
ExecStart=%h/.local/share/bin/tagion shell
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
